# files/etc/systemd/user/systembus-notify.service
[Unit]
Description=Show desktop notifications for earlyoom

[Service]
ExecStart=/usr/bin/systembus-notify

[Install]
WantedBy=default.target
